# debian/postinst
#!/bin/bash
# postinst script for cryptnox-cli
set -e

case "$1" in
    configure)
        # Enable and start pcscd service
        if command -v systemctl &>/dev/null; then
            systemctl enable pcscd 2>/dev/null || true
            systemctl start pcscd 2>/dev/null || true
        fi

        # Note: Python dependencies should be installed by the user
        # The deb package provides the core files, but Python deps
        # are better managed via pip to avoid conflicts
        echo ""
        echo "========================================"
        echo "  Cryptnox CLI installed"
        echo "========================================"
        echo ""
        echo "If you encounter import errors, install dependencies:"
        echo "  pip3 install --user cryptnox-cli"
        echo ""
        echo "Or use the installer script for full setup:"
        echo "  curl -fsSL https://raw.githubusercontent.com/cryptnox-snap/cryptnox-installer/main/install.sh -o install.sh"
        echo "  chmod +x install.sh && ./install.sh --native"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
